# apps/karaoke/com.sam.karaoke.redeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sam.karaoke.redeploy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/davidquinton/ReverseLab/SAM/apps/karaoke/redeploy.sh</string>
    </array>

    <!-- Run every Sunday at 9 AM (weekly refresh before 7-day expiry) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/davidquinton/ReverseLab/SAM/apps/karaoke/redeploy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davidquinton/ReverseLab/SAM/apps/karaoke/redeploy.log</string>

    <!-- Only run if devices are likely connected (user logged in) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
